# server/daemon/debian/postinst
#!/bin/sh
# postinst script for netidmd
#
# see: dh_installdeb(1)

set -e


case "$1" in
    configure)
        echo "Creating the netidmd group for config & cert ownership..."
        systemd-sysusers
        echo "Fixing ownership of server configuration ..."
        chown :netidmd /etc/netidmd/server.toml*

        echo "============================="
        echo "Thanks for installing Netidm!"
        echo "============================="
        echo "Please ensure you modify the configuration file at /etc/netidmd/server.toml"
        echo "Only then: systemctl enable netidmd.service"
        echo "Full examples are in /usr/share/netidmd/"
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
